# asst03/transcode.cron
* * * * * flock -x lock.lock -c 'minimoto_transcode'
